# systemd/sailfish-device-encryption-community-service.service
[Unit]
Description=Device Encryption Service

[Service]
Type=simple
ExecStart=/usr/bin/sailfish-device-encryption-community-service
